# Farmers-Dashboard-Backend/cron/crontab
# Environment Variables shared by all cron jobs
SHELL=/bin/bash

10 0 * * * cd /scheduler && ./run_scheduler.sh > /proc/1/fd/1 2>/proc/1/fd/2
