# Farmers-Dashboard-Backend/cron/crontab
10 0 * * * cd /scheduler && ./run_scheduler.sh > /proc/1/fd/1 2>/proc/1/fd/2
